# systemd/lid-monitor.path
[Unit]
Description=Monitor lid state

[Path]
PathChanged=/proc/acpi/button/lid/LID/state

[Install]
WantedBy=sleep.target
